# bin/schedulix-client
#!/bin/bash
#
### BEGIN INIT INFO
# Provides:          schedulix-client
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Should-Start:      $network $named $time
# Should-Stop:       $network $named $time
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start and stop the schedulix Enterprise Job Scheduling agents
# Description:       Controls the schedulix agent processes
### END INIT INFO
#
. /opt/schedulix/etc/SETTINGS
set -e
set -u

if [ -f /lib/lsb/init-functions ]; then
	. /lib/lsb/init-functions
	LOG_WARNING_MSG=log_warning_msg
	LOG_FAILURE_MSG=log_failure_msg
	ECHO_SUCCESS=echo
	ECHO_FAILURE=echo
else
	LOG_WARNING_MSG=echo
	LOG_FAILURE_MSG=echo
	ECHO_SUCCESS=echo
	ECHO_FAILURE=echo
fi

HOSTNAME=`hostname | tr [a-z] [A-Z]`


SELF=$(cd $(dirname $0); pwd -P)/$(basename $0)

# Safeguard (relative paths, core dumps..)
cd /
umask 077

case "${1:-''}" in
  'start')
	echo -n "Starting schedulix Enterprise Job Scheduling jobservers"
	if [ -d /opt/schedulix/taskfiles/$HOSTNAME ]; then
		cd /opt/schedulix/taskfiles/$HOSTNAME
		for USR in *; do
			if grep $USR /etc/passwd >/dev/null 2>&1; then
				cd $USR;
				for SVR in *; do
					su - $USR -c "
						[ -f \$HOME/.bashrc ] && . \$HOME/.bashrc;
						if [ -f \$HOME/etc/$SVR.init ]; then 
							. \$HOME/etc/$SVR.init;
							$BICSUITEHOME/bin/jobserver-run \$HOME/etc/$SVR.conf \$HOME/log/$SVR.out;
						fi
					"
				done
				cd ..;
			fi
		done
	fi
	$ECHO_SUCCESS
	;;

  'stop')
	echo -n "Stopping schedulix Enterprise Job Scheduling jobservers"
	if [ -d /opt/schedulix/taskfiles/$HOSTNAME ]; then
		cd /opt/schedulix/taskfiles/$HOSTNAME
		for USR in *; do
			if grep $USR /etc/passwd >/dev/null 2>&1; then
				cd $USR;
				for SVR in *; do
					PNAME="~$USR/log/$SVR.out"
					PNAME=`eval "ls $PNAME"`
					if [ -p $PNAME ]; then
						if fuser -s $PNAME >/dev/null 2>&1; then
							fuser -k $PNAME >/dev/null 2>&1;
						fi
					fi
				done
				cd ..;
			fi
		done
	fi
	$ECHO_SUCCESS
	;;

  'restart')
	echo "Restarting schedulix Enterprise Job Scheduling jobservers"
	/etc/init.d/schedulix-client stop
	/etc/init.d/schedulix-client start
	;;

  'reload'|'force-reload')
  	$LOG_WARNING_MSG "Reloading jobservers not supported; use 'restart' instead"
	log_end_msg 0
	;;

  'status')
	if [ -d /opt/schedulix/taskfiles/$HOSTNAME ]; then
		cd /opt/schedulix/taskfiles/$HOSTNAME
		for USR in *; do
			if grep $USR /etc/passwd >/dev/null 2>&1; then
				cd $USR;
				for SVR in *; do
					PNAME="~$USR/log/$SVR.out"
					PNAME=`eval "ls $PNAME"`
					if fuser -s $PNAME >/dev/null 2>&1; then
						echo "GLOBAL.$HOSTNAME.$USR.$SVR is up"
					else
						echo "GLOBAL.$HOSTNAME.$USR.$SVR is down"
					fi
				done
				cd ..;
			fi
		done
	else
		echo "No installed jobservers found"
	fi
  	;;

  *)
	echo "Usage: $SELF start|stop|restart|reload|force-reload|status"
	exit 1
	;;
esac
